# roles/devbootstrap/templates/etc/profile.d/dev_bootstrap.sh
alias dev-bootstrap='{{ dev_bootstrap.venv_dir }}/bin/python {{ dev_bootstrap.lib_dir }}/devbootstrap/devbootstrap.py'
