# roles/container_host/files/etc/systemd/system/containerd.service.d/10-ulimit.conf
# needed to prevent images from eating up tons of memory and getting
# oomkilled
#   https://github.com/containerd/containerd/pull/7566
[Service]
LimitNOFILE=1048576
